# .config/systemd/user/awesome-session.target
[Unit]
Description=AwesomeWM graphical session
Wants=graphical-session.target
After=graphical-session.target
